# image-processing.service.example
[Unit]
Description=Image Processing Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/image-processing
Environment="ALIBABA_ACCESS_KEY_ID=REPLACE_WITH_YOUR_ACCESS_KEY_ID"
Environment="ALIBABA_ACCESS_KEY_SECRET=REPLACE_WITH_YOUR_ACCESS_KEY_SECRET"
ExecStart=/usr/local/bin/python3 run.py
Restart=always

[Install]
WantedBy=multi-user.target 